# backend/backend.service
[Unit]
After=network.service

[Service]
ExecStart=/home/ubuntu/MedWasteWeb/backend/run_backend_service.sh

[Install]
WantedBy=default.target