# backend/backend.service
[Unit]
After=network.service

[Service]
ExecStart=/home/ubuntu/MedWasteWeb/backend/backend-app.py

[Install]
WantedBy=default.target